# pyscraper/crontab
* * * * * root /usr/src/app/run_cron_job.sh >> /var/log/cron.log 2>&1
